# com.geoclaw.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Runs scripts/backup.sh daily at 03:30 local time.
  Install:   ./install-backup.sh
  Uninstall: ./install-backup.sh --uninstall
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.geoclaw.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>__WORKDIR__/scripts/backup.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__WORKDIR__/.geoclaw-backup.log</string>

  <key>StandardErrorPath</key>
  <string>__WORKDIR__/.geoclaw-backup.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
